# openvpn.target
[Unit]
Description=Start openvpn tunnels
After=network.target
Wants=network.target

[Install]
WantedBy=multi-user.target
